# launchd/com.xiaolong.feishu-worklog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xiaolong.feishu-worklog</string>

    <!-- 每天 10:00 触发；如果电脑当时在睡，唤醒后 launchd 会自动补跑
         （RunAtLoad 只在 agent 加载/登录那一刻触发，不会在睡眠唤醒时再跑，
          「唤醒补跑」靠的是 StartCalendarInterval 的错过即补特性）。 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 用户登录 / agent 加载（含 reload）时也各跑一次。
         配合 --catch-up，已生成的日子会快速退出。手动触发也可：
           launchctl kickstart -k gui/$(id -u)/com.xiaolong.feishu-worklog -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 跑 FeishuWorklog.app 里的启动器（而非直接跑 python）：
         TCC 按 .app bundle 记权限，"完全磁盘访问"授权一次即永久有效，
         不受 Homebrew python 版本变化影响。.app 由 launchd/build_app.sh 生成。
         启动器自己会 chdir 到项目目录并跑 scripts.run_daily --catch-up。 -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/xiaolong/go/src/feishu-worklog/FeishuWorklog.app/Contents/MacOS/FeishuWorklog</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/xiaolong/go/src/feishu-worklog</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/xiaolong/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- 用户登录态相关：让 claude CLI 能读 keychain -->
        <key>HOME</key>
        <string>/Users/xiaolong</string>
    </dict>

    <!-- 日志（运行输出/错误） -->
    <key>StandardOutPath</key>
    <string>/Users/xiaolong/go/src/feishu-worklog/data/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xiaolong/go/src/feishu-worklog/data/launchd.err.log</string>

    <!-- 防止僵尸：超过 10 分钟没结束就 kill -->
    <key>ExitTimeOut</key>
    <integer>600</integer>

    <!-- 同一个任务实例只允许 1 个并发 -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
